# status-bar/src/nativeMain/resources/etc/xdg/autostart/lcarsde-status-bar.desktop
[Desktop Entry]
Type=Application
Name=LCARSDE status bar
Exec=/usr/bin/status-bar.kexe
OnlyShowIn=lcarsde;
