# menu/src/linuxX64Main/resources/etc/xdg/autostart/lcarsde-app-menu.desktop
[Desktop Entry]
Type=Application
Name=LCARSDE application menu
Exec=/usr/bin/lcarsde-menu.kexe
OnlyShowIn=lcarsde;
