# status-bar/src/linuxX64Main/resources/etc/xdg/autostart/lcarsde-status-bar.desktop
[Desktop Entry]
Type=Application
Name=LCARSDE status bar
Exec=/usr/bin/lcarsde-status-bar.kexe
OnlyShowIn=lcarsde;
